# src/configuration/user_exec/systemd/system/update-image-task.timer
[Unit]
Description=System update timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target