# src/configuration/user_exec/systemd/system/update-image-task.service
[Unit]
Description=System update

[Service]
Type=oneshot
ExecStart=/usr/libexec/update-system.sh